# service_files/xyz.openbmc_project.UsbDevice.service
[Unit]
Description=USB Device

[Service]
ExecStart=/usr/bin/usb-device
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target
