# service_files/xyz.openbmc_project.FruDevice.service
[Unit]
Description=Fru Device
After=dbus.service

[Service]
ExecStart=/usr/bin/fru-device
Type=simple

[Install]
WantedBy=multi-user.target
